# notifier/sound-notifier.service.template
[Unit]
Description=Sound Notifier
After=network-online.target

[Service]
ExecStart=[Notifier]/notifier/notifier.py [IP] 12100 [Secret] [Message]

[Install]
WantedBy=default.target
